# dbmeter.service
[Unit]
Description=Heat Pump Sound Level Meter Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=i2ssound2
WorkingDirectory=/home/i2ssound2/dbmeter
ExecStart=/usr/bin/python3 /home/i2ssound2/dbmeter/dbmeter_mqtt.py
Restart=always
RestartSec=30

[Install]
WantedBy=multi-user.target